# templates/autostart.conf.j2
task

start on startup

script
    chdir {{ deploy_app_dir }}
    {{ command }}
end script